# arch/archlaptop/home/.config/systemd/user/rawdog.service
[Unit]
Description=Rawdog Oneshot Updater

[Service]
Type=oneshot
ExecStart=/usr/bin/rawdog -Wuw

[Install]
WantedBy=default.target

