# arch/archlaptop/home/.config/systemd/user/rawdog.timer
[Unit]
Description=Rawdog Oneshot Updater Timer

[Timer]
Persistent=true
OnCalendar=hourly

[Install]
WantedBy=timers.target

